# credit_management/claude-credit-monitor.service
[Unit]
Description=Claude Code Credit Monitor
After=network.target

[Service]
Type=simple
User=%u
WorkingDirectory=/home/per/gitrepos/Tmux-Orchestrator
ExecStart=/home/per/gitrepos/Tmux-Orchestrator/credit_management/credit_monitor.py
Restart=always
RestartSec=60
Environment="CLAUDE_CREDIT_CHECK_INTERVAL=5"

# Logging
StandardOutput=append:/home/per/.claude/credit_monitor.log
StandardError=append:/home/per/.claude/credit_monitor.log

[Install]
WantedBy=default.target