# systemd/tmux-orchestrator-self-healing-monitor.service
[Unit]
Description=Tmux Orchestrator Self-Healing Monitor
After=network.target

[Service]
Type=simple
User=clauderun
WorkingDirectory=/home/clauderun/Tmux-Orchestrator
ExecStart=/home/clauderun/Tmux-Orchestrator/queue_error_monitor_modular_implementation_enhanced.sh
Restart=on-failure
RestartSec=60
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=TMUX_ORCHESTRATOR_HOME=/home/clauderun/Tmux-Orchestrator
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/home/clauderun/.local/bin

# Resource limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=default.target