# systemd/tmux-orchestrator-queue.service
[Unit]
Description=Tmux Orchestrator Queue Processor
After=multi-user.target
Conflicts=tmux-orchestrator-scheduler.service

[Service]
Type=simple
User=clauderun
Group=clauderun
WorkingDirectory=/home/clauderun/Tmux-Orchestrator
ExecStart=/usr/bin/python3 /home/clauderun/Tmux-Orchestrator/scheduler.py --daemon --mode queue
Restart=always
RestartSec=15
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal
KillMode=process
TimeoutStopSec=60

# Security hardening (relaxed for tmux access)
PrivateDevices=no
NoNewPrivileges=yes
ProtectSystem=false
ProtectHome=false

[Install]
WantedBy=multi-user.target